# scripts/systemd units/sn-zfs.timer
[Unit]
Description=Run ZFS Snapshots script daily

[Timer]
OnCalendar=05:03:00

[Install]
WantedBy=timers.target